# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/osint-labs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/osint-labs

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/osint-labs/nginx.log
stderr_logfile=/var/log/osint-labs/nginx.err.log
priority=10

[program:dashboard]
command=gunicorn --bind 127.0.0.1:5000 --workers 2 --threads 4 --timeout 120 app:app
directory=/opt/osint/dashboard
autostart=true
autorestart=true
stdout_logfile=/var/log/osint-labs/dashboard.log
stderr_logfile=/var/log/osint-labs/dashboard.err.log
priority=20

[program:ttyd]
command=ttyd -p 7681 -W -b /terminal bash
autostart=true
autorestart=true
stdout_logfile=/var/log/osint-labs/ttyd.log
stderr_logfile=/var/log/osint-labs/ttyd.err.log
priority=30

[program:spiderfoot]
command=python3 /opt/spiderfoot/sf.py -l 127.0.0.1:5001
directory=/opt/spiderfoot
autostart=true
autorestart=true
stdout_logfile=/var/log/osint-labs/spiderfoot.log
stderr_logfile=/var/log/osint-labs/spiderfoot.err.log
priority=40
startsecs=15

[program:social-analyzer]
command=node app.js --gui --docker
directory=/opt/social-analyzer
environment=PORT="9005"
autostart=true
autorestart=true
stdout_logfile=/var/log/osint-labs/social-analyzer.log
stderr_logfile=/var/log/osint-labs/social-analyzer.err.log
priority=50
startsecs=15
